# .config/systemd/user/reuters-sync.service
[Unit]
Description=reuters-sync

[Service]
User=rybalkin
Group=rybalkin
Type=oneshot
ExecStart=/bin/bash /usr/lib/addons/reuters-sync.bash
