# .config/systemd/user/syncema.service
[Unit]
Description=syncema

[Service]
Type=oneshot
ExecStart=/usr/bin/python /usr/bin/syncema
