# .config/systemd/user/syncema.timer
[Unit]
Description=syncema

[Timer]
OnCalendar=*-*-* *:*:10,40
Persistent=true

[Install]
WantedBy=timers.target
